# ffj-test/11s-test/etc/uci-defaults/91-set_system_opts
#!/bin/sh
uci set system.@system[0].timezone='Europe/Berlin'
uci commit system
